# automatisch_ppt/automatic_ppt.service
[Unit]
Description=Automatic PPT Startup Script
After=multi-user.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStartPre=/usr/bin/xhost +SI:localuser:pi
ExecStart=/etc/startup_scripts/automatic_ppt.sh
User=pi

[Install]
WantedBy=multi-user.target
